# seminar9_systemd/u04.service
[Unit]
Description=Logging service triggered by timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/s04.sh
