# seminar9_systemd/u04.timer
[Unit]
Description=Run u04.service every two minutes

[Timer]
OnUnitActiveSec=2min
OnBootSec=2min
Unit=u04.service

[Install]
WantedBy=timers.target
