# seminar9_systemd/u05a.service
[Unit]
Description=Immortal file watcher (polling)

[Service]
ExecStart=/usr/local/bin/s05a.sh
Restart=always

[Install]
WantedBy=multi-user.target
