# seminar9_systemd/u05b.service
[Unit]
Description=Immortal file restorer (path unit)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/s05b.sh

[Install]
WantedBy=multi-user.target
