# seminar9_systemd/u05b.service
[Unit]
Description=Restore immortal file

[Service]
Type=oneshot
ExecStart=/usr/local/bin/s05b.sh
