# seminar9_systemd/u05b.path
[Unit]
Description=Watch immortal file changes

[Path]
PathChanged=/var/log/immortal.txt
PathExists=/var/log/immortal.txt
Unit=u05b.service

[Install]
WantedBy=multi-user.target
